# config/includes.chroot/etc/systemd/system/multi-user.target.wants/systemd-set-time.service
[Unit]
Description=set time per timedatectl

[Service]
Type=forking
ExecStart=/usr/local/bin/set_time 

[Install]
WantedBy=multi-user.target
